# gpio/gpio.cron
@reboot /home/chip/chip_marquee/gpio/ttySO_state.sh
